# systemd/.config/systemd/user/udiskie.service
[Unit]
Description=udiskie daemon

[Service]
Type=simple
ExecStart=/usr/bin/udiskie

[Install]
WantedBy=default.target
